# crontab
0 * * * * /home/john/scripts/gdrivecron.sh > /home/john/scripts/glogs/glogs.log 2>&1
